# other/systemd/solarthing-request.service
[Unit]
Description=The service for SolarThing Request

[Service]
TimeoutStartSec=1
ExecStart=/opt/solarthing/program/request/run.sh
User=solarthing
Group=solarthing
StandardOutput=null

[Install]
WantedBy=multi-user.target
