# other/systemd/solarthing-pvoutput.service
[Unit]
Description=The service for SolarThing PVOutput

[Service]
TimeoutStartSec=1
ExecStart=/opt/solarthing/program/pvoutput/run.sh
User=solarthing
Group=solarthing

[Install]
WantedBy=multi-user.target
