# other/systemd/solarthing.service.template
[Unit]
Description=The service for SolarThing ##name##

[Service]
TimeoutStartSec=1
ExecStart=##solarthing_root##/program/##name##/run.sh
User=solarthing
Group=solarthing
StandardOutput=null

[Install]
WantedBy=multi-user.target
